# extras/wg-netns.target
[Unit]
Description=WireGuard Tunnels via wg-netns(8)
